# src/systemd/ai-agent.service
[Unit]
Description=AI Agent Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-agent
ExecStart=/usr/bin/node /opt/ai-agent/server.js
Restart=always
RestartSec=10

# Resource limits
MemoryMax=2G
CPUQuota=300%

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target